# roles/bitcoin/templates/bitcoin.service.j2
[Unit]
Description=Bitcoin service

[Service]
User=root
ExecStart={{ bitcoin_commands_dir }}/bitcoind
Restart=always

[Install]
WantedBy=multi-user.target
